# server/eventd/units/user/eventd-control.socket.in
[Unit]
Description=eventd control socket

[Socket]
Service=eventd.service
SocketMode=0600
FileDescriptorName=eventd-control
ListenStream=%t/@PACKAGE_NAME@/private

[Install]
WantedBy=sockets.target
