# server/eventd/units/user/eventd.socket.in
[Unit]
Description=eventd sockets

[Socket]
SocketMode=0660
ListenStream=%t/eventd/@EVP_UNIX_SOCKET@

[Install]
WantedBy=sockets.target
