# data/units/eventd.service.in
[Unit]
Description=eventd

[Service]
Sockets=eventd-control.socket eventd-evp.socket
ExecStart=@bindir@/eventd --listen-default --listen-default-unix
ExecReload=@bindir@/eventdctl reload

[Install]
WantedBy=default.target
Also=eventd-control.socket eventd-evp.socket
