# data/units/eventd-evp.socket.in
[Unit]
Description=eventd EvP Sockets

[Socket]
Service=eventd.service
ListenStream=%t/@PACKAGE_NAME@/@UNIX_SOCKET@
ListenStream=7100

[Install]
WantedBy=sockets.target
