# data/units/eventd-control.socket.in
[Unit]
Description=eventd Control Socket

[Socket]
Service=eventd.service
ListenStream=%t/@PACKAGE_NAME@/private

[Install]
WantedBy=sockets.target
