# server/eventd/units/system/eventd-control.socket.in
[Unit]
Description=eventd control socket

[Socket]
Service=eventd.service
SocketUser=eventd
SocketMode=0600
FileDescriptorName=eventd-control
ListenStream=/run/eventd/private

[Install]
WantedBy=sockets.target
