# server/eventd/units/system/eventd-control.socket.in
[Unit]
Description=eventd control socket

[Socket]
Service=eventd.service
SocketUser=eventd
SocketMode=0600
ListenStream=/run/@PACKAGE_NAME@/private

[Install]
WantedBy=sockets.target
