# server/eventd/units/system/eventd.socket.in
[Unit]
Description=eventd sockets

[Socket]
SocketUser=eventd
SocketGroup=eventd
SocketMode=0660
FileDescriptorName=evp
ListenStream=/run/eventd/@EVP_UNIX_SOCKET@

[Install]
WantedBy=sockets.target
