# server/eventd/units/system/eventd.service.in
[Unit]
Description=eventd
Documentation=man:eventd(1)
Documentation=man:eventd.conf(5)
Documentation=https://www.eventd.org/

[Service]
Type=notify
User=eventd
Sockets=eventd-control.socket eventd.socket
ExecStart=@bindir@/eventd --listen systemd
ExecReload=@bindir@/eventdctl reload

[Install]
Also=eventd-control.socket
Also=eventd.socket
