# server/eventd/units/eventd.socket.in
[Unit]
Description=eventd sockets

[Socket]
SocketMode=0660
ListenStream=%t/@PACKAGE_NAME@/@EVP_UNIX_SOCKET@

[Install]
WantedBy=sockets.target
