# server/eventd/units/eventd.service.in
[Unit]
Description=eventd
Documentation=man:eventd(1)
Documentation=man:eventd.conf(5)
Documentation=https://www.eventd.org/

[Service]
Type=notify
ExecStart=@bindir@/eventd --listen systemd
ExecReload=@bindir@/eventdctl reload

[Install]
WantedBy=default.target
Also=eventd.socket
