# server/eventd/units/eventd.service.in
[Unit]
Description=eventd

[Service]
Type=notify
Sockets=eventd-control.socket eventd-evp.socket
ExecStart=@bindir@/eventd --listen systemd
ExecReload=@bindir@/eventdctl reload

[Install]
WantedBy=default.target
Also=eventd-control.socket eventd-evp.socket
