# server/eventd/units/eventd.service.in
[Unit]
Description=eventd

[Service]
Type=notify
ExecStart=@bindir@/eventd --listen systemd
ExecReload=@bindir@/eventdctl reload

[Install]
WantedBy=default.target
Also=eventd.socket
